# runit/ays/run
#!/bin/bash

HOME=/root exec python ./run.py
